# scripts/debian/bitcore-bitcoind.conf
description "Bitcoin Core for Bitcore"
author "BitPay, Inc."

limit nofile 20000 30000

start on runlevel [2345]
stop on runlevel [016]

kill timeout 300
kill signal SIGINT

# user/group for bitcore daemon to run as
setuid bitcore
setgid bitcore

# home dir of the bitcore daemon user
env HOME=/home/bitcore

respawn
respawn limit 5 15

script
  exec bitcore-bitcored -datadir=/home/bitcore/.bitcore/data/
end script
