# etc/init/bitcored.conf
# Bitcored Upstart script

description "Bitcore Daemon"
author "BitPay, Inc."

limit nofile 20000 30000

start on runlevel [2345]
stop on runlevel [016]

kill timeout 300
kill signal SIGINT

# user/group for bitcore daemon to run as
setuid bitcore
setgid bitcore

# home dir of the bitcore daemon user
env HOME=/home/bitcore

# path should include the correct version of node
env PATH=/home/bitcore/.nvm/versions/node/v0.12.7/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

respawn
respawn limit 5 15

script
  exec bitcored
end script
